# rakuos-settings/debian/postinst
#!/bin/sh

set -e

systemctl enable --now pci-latency.service
systemctl restart rtkit-daemon.service

udevadm control --reload-rules && udevadm trigger
