# rakuos-settings/debian/postinst
#!/bin/sh

set -e

systemctl enable pci-latency.service
