# rankuos-gaming-meta/debian/postinst
#!/bin/bash
set -e
echo "🎮 RankuOS Gaming Environment installed!"
systemctl --user enable --now gamemoded.service || true
