# ananicy-cpp/debian/postinst
#!/bin/sh

set -e

systemctl enable --now ananicy-cpp.service
