# ananicy-cpp/debian/postinst
#!/bin/sh

set -e

systemctl enable ananicy-cpp.service
